# etc/ledmatrix.service
[Unit]
Description=ledmatrix

[Service]
Type=simple
ExecStart=/usr/local/bin/ledmatrix

[Install]
WantedBy=sysinit.target
